# recipe.conf
start on runlevel [23]
console log

respawn
respawn limit 100 10

exec su ubuntu -c "sleep 10; cd /home/ubuntu/recipeserver; ./RunServer.bash || true"
